# crontab
*/5 * * * * /app/SigSciSumo.sh > /var/log/sigsci.log 2>&1
# You need an empty line at the end of the file so that it is considered a valid cron file.

